# systemd/lib/systemd/system/testpooldb.service
[Unit]
Description=testpooldb

[Service]
Type=simple
ExecStart=/usr/bin/testpooldb runserver
Restart=always

[Install]
WantedBy=multi-target.target
