# debian/python-testpool.postinst
#!/bin/bash
echo "info: python-testpool.postinst"
init_system=`dpkg -S /sbin/init`
echo "info: init system $init_system"
##
# Install on systemd. 
# setup.py installs /usr/local/share/testpooldb.service
# Without checking for now install into systemd area
cp /usr/testpool/systemd/testpooldb.service  /lib/systemd/sytem/
systemctl daemon-reload
systemctl enable testpooldb.service
##
