# debian/testpool.postinst
#!/bin/bash

systemctl deamon-reload
systemctl enable testpooldb.service
systemctl start testpooldb.service
