# usr/lib/systemd/system/testpooldb.service
[Unit]
Description=testpooldb

[Service]
ExecStart=/usr/bin/testpooldb runserver
